# deploy/vps/tiktok-review-bot.service.template
[Unit]
Description=TikTok Review Notifier Bot
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/python __APP_DIR__/tiktok_review_notifier.py --headless
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
